# ops/meloan.service
[Unit]
Description=Meloan PWA
After=network.target postgresql.service

[Service]
Type=simple
WorkingDirectory=/var/www/meloan
EnvironmentFile=/var/www/meloan/.env.production
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=5
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
